# src/backend/proxy/conf/supervisor.conf

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket


[program:scraper]
command=python /opt/scraper/python/app.py
autostart=true
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/var/log/scraper.log
